# core0/us/etc/cron.daily/install_blobs.in
#!/bin/sh
/usr/local/bin/install_blob >/var/log/##str3##/install_blobs.last.log 2>&1
exit $?
