# core0/us/etc/cron.daily/us-ports.in
#!/bin/sh
##str3##-ports $@  >/var/log/##str3##/##str3##-ports.last.log 2>&1
exit $?
