# core0/us/etc/systemd/system/us-gov.service.in
[Unit]
Description = ##str3## governance
After = network-online.target

[Service]
Type=simple
User=##str11##
EnvironmentFile=/home/##str11##/.##str3##/config
ExecStart=
ExecStart=/usr/local/bin/##str3##-gov -ds -e ${GOV_MAXEDGES_PEER} -v ${GOV_MAXEDGES_DEVICE} -c ${CHANNEL} -p ${PORT} -pp ${PPORT} -status ${GOV_STATUS_FILE} -E ${SAVE_EVIDENCES}
Restart=always
RestartSec=5
Nice=-1
WorkingDirectory=/home/##str11##

[Install]
WantedBy = multi-user.target

