# core0/us/etc/systemd/system/us-wallet.service.in
[Unit]
Description = ##str3## wallet
After = ##str3##-gov.service

[Service]
Type=simple
EnvironmentFile=/home/##str11##/.##str3##/config
ExecStart=
ExecStart=/usr/local/bin/##str3##-wallet -d -c ${CHANNEL} -lp ${WALLET_PORT} -pp ${WALLET_PPORT} -bp ${PORT} -e ${WALLET_MAXEDGES_PEER} -v ${WALLET_MAXEDGES_DEVICE}
User=##str11##
Restart=always
RestartSec=5
Nice=0
WorkingDirectory=/home/##str11##

[Install]
WantedBy = multi-user.target

